# etc/rc.local
#!/bin/sh -e

rfkill block bluetooth
rfkill block wlan
